# supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod --config %(ENV_MONGODB_CONF)s